# deploy/systemd/anneal-worker.service
[Unit]
Description=Anneal Worker
After=network-online.target postgresql.service anneal-api.service
Wants=network-online.target

[Service]
Type=simple
User=anneal
Group=anneal
EnvironmentFile=/etc/anneal/anneal.env
WorkingDirectory=/opt/anneal
ExecStart=/opt/anneal/bin/worker
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
