# deploy/docker/prebuilt/node-supervisord.conf
[program:xray]
command=/usr/local/bin/xray run -c /var/lib/anneal/xray/config.json
autorestart=true
priority=10
startsecs=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:singbox]
command=/usr/local/bin/hiddify-core run -c /var/lib/anneal/singbox/config.json
autorestart=true
priority=20
startsecs=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:node-agent]
command=/usr/local/bin/node-agent
autorestart=true
priority=30
startsecs=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
